# crontab
# Main data pipeline runs daily at 2 AM
0 2 * * * python /app/main_pipeline.py >> /var/log/cron.log 2>&1

# Guardian Agent runs weekly on Sunday at 5 AM
0 5 * * 0 python /app/agent_tasks.py >> /var/log/cron.log 2>&1

# Quality Auditor runs twice a week at 4 AM
0 4 * * 2,5 python /app/quality_auditor_agent.py >> /var/log/cron.log 2>&1
